# etc/cron.daily/daily-db-backup
#!/bin/sh
/opt/konveyor/perform-backup.sh
exit 0