# etc/cron.d/konveyor-db-backup
0 */4 * * * root /opt/konveyor/perform-backup.sh